# deploy/ai-comic-server.service
[Unit]
Description=AI Comic Generator Backend (Go/Gin)
After=network.target mysqld.service redis.service
Wants=mysqld.service redis.service

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/opt/ai-comic
ExecStart=/opt/ai-comic/bin/server

# 敏感配置也可用环境变量覆盖 config.yaml
# EnvironmentFile=-/opt/ai-comic/.env

Restart=always
RestartSec=3
# 崩溃循环保护：10s 内最多重启 5 次，超过则停下等人工介入
StartLimitIntervalSec=10
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-comic-server

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
